# per-once/update-idw-cfg.cron
* * * * * root /opt/irods-aws/update-idw-cfg.sh
