# per-once/update-idw-cfg.cron
* * * * * root /opt/irods-aws/per-mimute/update-idw-cfg.sh
